# home/dot_config/autostart/wezterm-mux-server.desktop
[Desktop Entry]
Exec=wezterm-mux-server --daemonize
Icon=dialog-scripts
Name=wezterm-mux-server
Type=Application
X-KDE-AutostartScript=true
