# dot_config/autostart/gammastep.desktop
[Desktop Entry]
Exec=/bin/gammastep -m randr
Icon=dialog-scripts
Name=gammastep
Type=Application
X-KDE-AutostartScript=true
