# dot_config/autostart/nm-applet.desktop
[Desktop Entry]
Exec=/bin/nm-applet
Icon=dialog-scripts
Name=nm-applet
Type=Application
X-KDE-AutostartScript=true
